# content/etc/udev/rules.d/z99-rpigps-disconnect.rules
ACTION=="remove", SUBSYSTEM=="tty" RUN+="/usr/bin/rpigps-system-halt.sh"
